# gmail_monitor.service
Description=Monitor gmail for automation commands
After=network.target

[Service]
Type=simple
ExecStart=/home/ubuntu/gmail_monitor/gmail_monitor.sh
Restart=always
RestartSec=30
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
